# auth_monitor.template.service
[Unit]
Description=Authentication Log Monitor Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 __DIR__/auth_monitor.py
WorkingDirectory=__DIR__
Restart=always

[Install]
WantedBy=multi-user.target
